# systemd/sentient-perception-http.service
[Unit]
Description=Sentient Core - Perception HTTP Service
After=network.target mosquitto.service redis-server.service
Wants=mosquitto.service redis-server.service

[Service]
Type=simple
User=cortana
WorkingDirectory=/opt/sentient-core
ExecStart=/usr/bin/python3 /opt/sentient-core/services/perception_http.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
