# systemd/sentient-contemplation.service
[Unit]
Description=Sentient Core - Contemplation Service (v2)
After=sentient-core.target ollama.service
PartOf=sentient-core.target

[Service]
Type=simple
User=cortana
Group=cortana
WorkingDirectory=/opt/sentient-core
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/sentient-core
EnvironmentFile=-/opt/sentient-core/config/sentient.env
ExecStart=/usr/bin/python3 -m sentient.services.contemplation.api
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=5
MemoryMax=4G
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentient-contemplation

[Install]
WantedBy=sentient-core.target
