# systemd/sentient-proactive.service
[Unit]
Description=Sentient Core - Proactive Behavior Engine (v2)
After=sentient-contemplation.service
Wants=sentient-contemplation.service
PartOf=sentient-core.target

[Service]
Type=simple
User=cortana
Group=cortana
WorkingDirectory=/opt/sentient-core
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/sentient-core
EnvironmentFile=-/opt/sentient-core/config/sentient.env
ExecStart=/usr/bin/python3 -m sentient.services.proactive
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentient-proactive

[Install]
WantedBy=sentient-core.target
