# systemd/sentient-conversation.service
[Unit]
Description=Sentient Core - Conversation Orchestrator (v2)
After=sentient-memory.service sentient-contemplation.service sentient-perception.service
Wants=sentient-memory.service sentient-contemplation.service sentient-perception.service
PartOf=sentient-core.target

[Service]
Type=simple
User=cortana
Group=cortana
WorkingDirectory=/opt/sentient-core
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/sentient-core
EnvironmentFile=-/opt/sentient-core/config/sentient.env
ExecStart=/usr/bin/python3 -m sentient.services.conversation
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentient-conversation

[Install]
WantedBy=sentient-core.target
