# systemd/sentient-voice.service
[Unit]
Description=Sentient Core - Voice Pipeline (v2)
After=sentient-wake-word.service sentient-conversation.service
Wants=sentient-wake-word.service sentient-conversation.service
PartOf=sentient-core.target

[Service]
Type=simple
User=cortana
Group=cortana
WorkingDirectory=/opt/sentient-core
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/sentient-core
EnvironmentFile=-/opt/sentient-core/config/sentient.env
ExecStart=/usr/bin/python3 -m sentient.services.voice.pipeline
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentient-voice

[Install]
WantedBy=sentient-core.target
