# systemd/sentient-avatar.service
[Unit]
Description=Sentient Core - Avatar Bridge (v2)
After=sentient-core.target
PartOf=sentient-core.target

[Service]
Type=simple
User=cortana
Group=cortana
WorkingDirectory=/opt/sentient-core
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/sentient-core
EnvironmentFile=-/opt/sentient-core/config/sentient.env
ExecStart=/usr/bin/python3 -m sentient.services.avatar.bridge
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentient-avatar

[Install]
WantedBy=sentient-core.target
